# worldforge.service
[Unit]
Description=WorldForge DnD Lore System
After=syslog.target network.target

[Service]
Type=simple
WorkingDirectory=/srv/worldforge
ExecStart=/usr/bin/make prod
ExecStop=/usr/bin/make down
Restart=always

[Install]
WantedBy=multi-user.target